# files/apcupsd_exporter.service.j2
[Unit]
Description=apcupsd_exporter

[Service]
TimeoutStartSec=0
ExecStart=/usr/local/bin/apcupsd_exporter -telemetry.addr :{{ apcupsd_exporter_http }}

[Install]
WantedBy=multi-user.target
